# .config/systemd/user/dynamicwall.service
[Unit]
Description=dynamic-wall unit

[Service]
Type=oneshot
ExecStart=dynamicwall

[Install]
WantedBy=multi-user.target
