# .config/systemd/user/gebaard-watcher.service
[Unit]
Description=Gebaar restarter

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user try-restart gebaard.service
