# .config/systemd/user/gebaard-watcher.path
[Unit]
Description=Watch for gebaar configuration changes

[Path]
PathModified=%h/.config/gebaar/gebaard.toml

[Install]
WantedBy=graphical-session.target
