# .config/systemd/user/gebaard-watcher.path
# Watches for file change and triggers gebaard-watcher.service
[Path]
PathModified=%h/.config/gebaar/gebaard.toml

[Install]
WantedBy=multi-user.target
