# .config/systemd/user/gebaard-watcher.path
[Path]
PathChanged=%h/.config/gebaar/gebaard.toml

[Install]
WantedBy=default.target
