# .config/autostart/Tapping.desktop
[Desktop Entry]
Comment[en_US]=
Comment=
Exec=xinput --set-prop "Synaptics TM3044-002" "libinput Tapping Enabled" 1
GenericName[en_US]=
GenericName=
Icon=system-run
MimeType=
Name[en_US]=
Name=
Path=
StartupNotify=true
Terminal=false
TerminalOptions=
Type=Application
X-DBUS-ServiceName=
X-DBUS-StartupType=
X-KDE-RunOnDiscreteGpu=false
X-KDE-SubstituteUID=false
X-KDE-Username=
