# Scripts/System/etc/profile.d/kf5.sh
# Stop kscreen.xcb.helper message from spamming logs
export QT_LOGGING_RULES="kscreen.xcb.helper=false"
